# php-nginx/5.6/config/nginx.ini
[program:nginx]
command     = /usr/sbin/nginx -g 'daemon off;'
user        = root
autostart   = true
startsecs   = 5
autorestart = true